# crontab
0 3 * * 6 /webapps/kmailbackend/kmail/cronjob
